# .config/containers/systemd/signal.container
[Unit]
Description=Netbird Signal Service

[Container]
ContainerName=signal
Image=docker.io/netbirdio/signal:latest
Network=netbird.network
LogDriver=journald
AutoUpdate=registry

[Install]
WantedBy=multi-user.target default.target
